# deploy/systemd/reds-scheduler.timer
[Unit]
Description=Arm Reds post-game timers (09:00 & 16:00 ET)

[Timer]
OnCalendar=*-*-* 09,16:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
